# compliance/puppet/stig/rhel-6/modules/aide/files/aidereset.cron
#!/bin/sh

/usr/sbin/aide --update
cp /var/lib/aide/aide.db.new.gz /var/lib/aide/aide.db.gz

exit 0
