# compliance/puppet/stig/rhel-6/modules/aide/files/aide.cron
#!/bin/sh

/usr/sbin/aide --check

exit 0
